# deploy/systemd/orcheo-beat.service
[Unit]
Description=Orcheo Celery Beat Scheduler
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=orcheo
Group=orcheo
WorkingDirectory=/opt/orcheo
Environment=PATH=/opt/orcheo/.venv/bin:/usr/bin:/bin
Environment=REDIS_URL=redis://localhost:6379/0
Environment=ORCHEO_MULTI_WORKSPACE_ENABLED=false
Environment=ORCHEO_MULTI_WORKSPACE_DEFAULT_WORKSPACE_SLUG=default
Environment=ORCHEO_MULTI_WORKSPACE_WORKSPACE_HEADER=X-Orcheo-Workspace
Environment=ORCHEO_WORKSPACE_BACKEND=postgres
EnvironmentFile=-/etc/orcheo/orcheo.env

ExecStart=/opt/orcheo/.venv/bin/celery \
    -A orcheo_backend.worker.celery_app \
    beat \
    --loglevel=info

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
